# systemd/vedurstod-turbo-deploy.service
[Unit]
Description=Veðurstöð Turbo auto-deploy (pull + build + restart)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/vedurstod-turbo
ExecStart=%h/vedurstod-turbo/scripts/deploy-pi.sh
StandardOutput=journal
StandardError=journal
